# mynas_build/config/includes.chroot/etc/rc.local
#!/bin/bash
cd /app && /app/mynas &
echo "init done"
exit 0
