# mynas_build/config/includes.chroot/etc/rc.local
#!/bin/bash
cd /app && /app/mynas_go &
cd /app/scripts && ./auto_start.sh
echo "init done"
exit 0
