# mynas_build/cron/watch_kernel.cron
0 2 * * * /data/tools/watch_kernel.sh 6.15 compile_kernel
