# mynas_build/cron/watch_kernel.cron
0 2 * * * /data/tools/watch_kernel.sh 6.17 compile_kernel &
